# deploy/linux/qrm.service
[Unit]
Description=QRM Rack UI and Backend
After=network.target

[Service]
Type=simple
User=bert.olsson
Group=bert.olsson
WorkingDirectory=/home/bert.olsson/src/qrm-ui
ExecStart=/bin/bash /home/bert.olsson/src/qrm-ui/deploy/linux/start.sh
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=qrm-service

[Install]
WantedBy=multi-user.target
